# launchd/com.ringkas.prd-enrich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ringkas.prd-enrich</string>
  <!-- Runs after A's sync (11:00, ~40 min). Absolute Volta npm path; real vault.
       MiniMax router, model minimax/MiniMax-M3, with stream:false handled in code. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string><string>-lc</string>
    <string>cd /Users/edwardpham/Documents/Workspace/Ringkas/Programming/Personal/llm-wiki &amp;&amp; VAULT_PATH="/Users/edwardpham/Documents/Backup/Obsidian/ringkas" LLM_BASE_URL="https://9router-1.dat-nguyen.me/v1" LLM_MODEL="minimax/MiniMax-M3" /Users/edwardpham/.volta/bin/npm run enrich</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>45</integer></dict>
  <key>StandardOutPath</key><string>/tmp/prd-enrich.log</string>
  <key>StandardErrorPath</key><string>/tmp/prd-enrich.err.log</string>
</dict>
</plist>
